# templates/germqtt.service.erb
[Unit]
Description=Germqtt Daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/germqtt /etc/germqtt.conf --foreground
User=germqtt
Group=germqtt

[Install]
WantedBy=multi-user.target
